# nzbwatch.service
[Unit]
Description=NZB Watch Folder Daemon

[Service]
Type=simple
ExecStart=/usr/bin/ruby /usr/bin/nzbwatch.rb

[Install]
WantedBy=default.target
